# systemd/kokpit-runs-archive.timer
[Unit]
Description=Kokpit runs/ aylık arşivleme (30+ gün eski log'lar)

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
